# darknet.ini
[program:DarkNet]
command=/root/anaconda3/bin/python3 /root/darknet/run.py
autorstart=true
redirect_stderr=true
stdout_logfile=/mnt/log/darknet.log
directory=/root/darknet
user=root
priority=1
